# note/hbase不睡觉书/1/自启动项/hadoop-namenode
#!/bin/bash  
#chkconfig: 2345 82 08
#description: hadoop namenode service

RETVAL=0
# Source hadoop environment service
. /home/hadoop/.bashrc
start(){
    su hadoop -c "$HADOOP_PREFIX/sbin/hadoop-daemon.sh start namenode"
}
stop(){
    su hadoop -c "$HADOOP_PREFIX/sbin/hadoop-daemon.sh stop namenode"
}

case $1 in
start)
    start
;;
stop)
    stop
;;
esac
exit $RETVAL

